# packages/qmp-meshtodon/files/etc/uci-defaults/zmeshtodon_banner
#!/bin/sh
cat /etc/meshtodon/meshtodon.banner >> /etc/banner
exit 0
